# external/renaWebsocket.conf
[program:renaWebSocket]
command=/usr/bin/php /storage/www/projectRena/Rena run:websocket
directory=/storage/www/projectRena
autostart=true
autorestart=true
user=www-data
group=www-data
